# pallet/rc.local
#!/bin/sh -e
#
# rc.local
#
# Make sure that the script will "exit 0" on success or any other
# value on error.

nohup Xvfb :99

exit 0
